# mcp/systemd/aleph-docs-mcp-lint.service
[Unit]
Description=MCP Aleph Docs memory lint (oneshot)
After=network.target

[Service]
Type=oneshot
User=www-data
WorkingDirectory=/opt/mcp
EnvironmentFile=/opt/mcp/.env
ExecStart=/opt/mcp/.venv/bin/python -m memory.lint_cli --mode auto
